# nsncd.socket
[Unit]
PartOf=nsncd.service

[Socket]
ListenStream=/run/nscd/socket
SocketMode=0777

[Install]
WantedBy=sockets.target
